# monitor/src/listend.in
#!/bin/sh
#
# chkconfig: 12345 99 01
# description: Starts, stops, restarts and reloads tacc_stats listen daemon

# Source function library
#. /etc/init.d/functions

PATH=/bin:/sbin:/usr/bin
PID_FILE=listend_lock
SERVER=server
ARCHIVE=archive_dir

status() {
    pid_status
}

pid_status() {
    if pidof listend > /dev/null; then
	if [ "$1" != "quiet" ] ; then
	    echo "listend is running"
	fi
	RETVAL=1
    else
	if [ "$1" != "quiet" ] ; then 
	    echo "listend is not running"	
	fi
	RETVAL=0		   
    fi
}

start() {
    pid_status quiet
    if [ $RETVAL = 0 ] ; then 
	echo "starting listend"
	amqp_listend -s ${SERVER} -a ${ARCHIVE}
    fi
}
	
stop() {
    pid_status quiet
    if [ $RETVAL = 1 ] ; then 
	echo "stopping listend"
        cat ${PID_FILE} | kill -9 `awk '{print $1}'`	
       	rm -f ${PID_FILE}
    else
	echo "no instance of listend running"
    fi
}    

restart() {
    stop
    start
}

case "$1" in 
    start)
	start
	;;
    stop)
	stop
	;;
    status)
	status
	;;
    restart)
	restart
	;;
    *)
	echo $"Usage: $0 {start|stop|status|restart}"
esac

exit $?
